# deploy/cnip-api.service
[Unit]
Description=cnip.io API backend (127.0.0.1:18083)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/cnip/api
EnvironmentFile=/opt/cnip/.env
Environment=APP_HOST=127.0.0.1
Environment=APP_PORT=18083
ExecStart=/opt/cnip/api/bin/ip2region-backend
Restart=always
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
