# deploy/ip2region-backend.service
[Unit]
Description=cnip.io / ip2region.io Go backend (127.0.0.1:18083)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ip2region.io/api
EnvironmentFile=/opt/ip2region.io/.env
Environment=APP_HOST=127.0.0.1
Environment=APP_PORT=18083
ExecStart=/opt/ip2region.io/api/bin/ip2region-backend
Restart=always
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
